# debian/postinst
#!/bin/sh
systemctl enable node-dummy-app
systemctl start node-dummy-app
